# deploy/supervisor/wellguard.conf
[program:wellguard]
command=/opt/wellguard/venv/bin/streamlit run streamlit_app.py --server.port=8501 --server.headless=true --server.address=0.0.0.0
directory=/opt/wellguard
user=root
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stderr_logfile=/var/log/wellguard/err.log
stdout_logfile=/var/log/wellguard/out.log
environment=HOME="/root",USER="root"
